# postgres/backup/crontabs
0,30 * * * * ash /bin/backup.sh
0 0 * * * ash /bin/cleanup.sh